# experiments/ede_exp/http_rest_chaos/jupyter_aspide.service
[Unit]
Description=Jupyter Aspide Notebook

[Service]
Type=simple
PIDFile=/run/jupyter_aspide.pid
# Step 1 and Step 2 details are here..
# ------------------------------------
ExecStart=/home/gabriel/anaconda3/envs/pyaspide/bin/jupyter-notebook --config=/home/gabriel/.jupyter/jupyter_notebook_config.py
User=gabriel
Group=gabriel
WorkingDirectory=/home/gabriel/Research/Aspide/workspace/v2/experiments/ede_exp/notebooks
Restart=always
RestartSec=10
#KillMode=mixed

[Install]
WantedBy=multi-user.target